# build/ansible/roles/consul/templates/consul-leave.conf.j2
description "Gracefully leave consul cluster"
# {{ ansible_managed }}
start on stopping consul

task

script
  consul leave
end script
